# extensions_data/autostart/webpage-accessor.desktop
[Desktop Entry]
Name=webpage_accessor
Comment=webpage accessor autostart script
Type=Application
Exec=webpage_accessor