# extensions_data/autostart/webpage-accessor.desktop
[Desktop Entry]
Name=webpage-accessor
Comment=webpage accessor autostart script
Type=Application
Exec=webpage-accessor